# orderbot.timer
sudo tee /etc/systemd/system/orderbot.timer >/dev/null <<'UNIT'
[Unit]
Description=Run Telegram Order Notifier daily at 22:30 Europe/Rome

[Timer]
# Run every day at 22:30 local time (server timezone should be Europe/Rome)
OnCalendar=*-*-* 22:30:00
Persistent=true
Unit=orderbot.service

[Install]
WantedBy=timers.target
UNIT